# svc/new-year.timer
[Unit]
Description=start of New Years Day

[Timer]
OnCalendar=yearly
Persistent=true

[Install]
WantedBy=timers.target
