# svc/update-copyright.service
[Unit]
Description=Update site copyright year

[Service]
Type=oneshot
ExecStart=/usr/bin/update-copyright
