# svc/update-copyright.timer
[Unit]
Description=Yearly Copyright Notice Update

[Timer]
OnCalendar=yearly
Persistent=true

[Install]
WantedBy=timers.target
